# scripts/ai.llmfit.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.llmfit.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rajamans/workspace/dev/repos/personal/llmfit/scripts/llmfit-server-managed.sh</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/rajamans/workspace/logs/llmfit-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rajamans/workspace/logs/llmfit-server.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/rajamans/workspace/dev/repos/personal/llmfit</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>